# casesquirrels/uwsgi.conf
[program:uwsgi]
command=uwsgi --emperor /etc/uwsgi/vassals --uid www-data --gid www-data
user=root
autostart=true
autorestart=true

